# deploy/serenzer-feedback-api.service.example
[Unit]
Description=Serenzer Feedback API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/serenzerfeedback
ExecStart=/usr/bin/python3 /var/www/serenzerfeedback/server.py
Restart=always
RestartSec=3
Environment=SERENZER_HOST=127.0.0.1
Environment=SERENZER_PORT=8000
Environment=SERENZER_IMPORT_API_KEY=replace-with-shared-secret

[Install]
WantedBy=multi-user.target
